# roles/systemd_mount_remote_nfs_storage/templates/nfs.mount.j2
[Unit]
Description=Mount NFS share {{ nfs_export }} from {{ nfs_server }} at {{ mount_path }}

[Mount]
What={{ nfs_server }}:{{ nfs_export }}
Where={{ mount_path }}
Type=nfs
Options=defaults,noatime
